# killmailer/killmailer/mailer.gunicorn.socket
[Unit]
Description=gunicorn socket

[Socket]
ListenStream=/apps/killmailer/app/killmailer/run/gunicorn.sock

[Install]
WantedBy=sockets.target
